# .config/systemd/user/upd.service
[Unit]
Description=Checking for pacman updates

[Service]
Type=oneshot
ExecStart=-%h/.local/bin/upd
