# .config/systemd/user/wallpaper.timer
[Unit]
Description=Random wallpaper with feh

[Timer]
OnStartupSec=5m
OnUnitActiveSec=30m
Unit=wallpaper.service

[Install]
WantedBy=timers.target
